# launchd/com.merlin.stack.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.merlin.stack.plist
  Starts the laptop-safe core profile 30 seconds after login.
  The 30s delay gives Docker Desktop time to fully initialize.

  Install:
    cp ~/merlin-ai/launchd/com.merlin.stack.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.merlin.stack.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.merlin.stack</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <!-- Wait 30s for Docker to be ready, then start only the core profile through wizard -->
    <string>sleep 30 &amp;&amp; cd "$HOME/merlin-ai" &amp;&amp; HOME_AI_PROFILE=core /bin/bash cli/wizard start core &gt;&gt; /tmp/merlin-stack.log 2&gt;&amp;1</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <!-- Restart the watcher if it crashes, but not the compose command itself -->
  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/merlin-stack.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/merlin-stack.log</string>

  <!-- Working directory -->
  <key>WorkingDirectory</key>
  <string>/Users/Shared</string>

  <!-- Environment — ensure PATH includes Homebrew and Docker CLI -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
